# _linux/rc.local
touch /var/lock/subsys/local
sh '/etc/rc.d/launch_startup.sh' &
exit 0